# gestion_usuarios.path
[Unit]
        Description=Alta Baja Usuarios Path
        Before=administracionUsuarios.service

[Path]
        PathExistsGlob=/srv/administracioUsuaris/*.csv

[Install]
        WantedBy=multi-user.target
